# inittab
# /etc/inittab

::sysinit:/sbin/rc sysinit
::wait:/sbin/rc boot
::wait:/sbin/rc default

# Set up a couple of getty's
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4

# Put a getty on the serial port
# You need to include this file in the initramfs image to activate it
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/rc shutdown
null::shutdown:/bin/umount -a -r 2>/dev/null
null::shutdown:/sbin/swapoff -a 2>/dev/null

