# init.d/cron
#!/sbin/runscript

start() {
	ebegin "Starting cron"
	crond $CRON_OPTS
	eend $?
}

stop () {
	ebegin "Stopping cron"
	killall crond
	eend $?
}
