# init.d/dnsd
#!/sbin/runscript

start() {
	ebegin "Starting dnsd"
	dnsd -d $DNSD_OPTS
	eend $?
}

stop () {
	ebegin "Stopping dnsd"
	killall dnsd
	eend $?
}
