# init.d/hostname
#!/sbin/runscript

start() {
	if [ -f /etc/hostname ] ; then
		ebegin "Setting hostname"
		hostname -F /etc/hostname
		eend $?
	fi
}
