# init.d/hwdrivers
#!/sbin/runscript

# Load hardware drivers
start() {
	# check for boot option "nocoldplug"
	for i in `cat /proc/cmdline`; do
		case $i in
			modules=*)
				MODULES="$MODULES `echo ${i#modules=} | tr ',' ' '`";;
			noautodetect)
				AUTODETECT=no;;
		esac
	done
	
	if [ "$MODULES" ]; then
		ebegin "Loading boot drivers"
		modprobe -a $MODULES 2> /dev/null
		eend $?
	fi

	if [ "$AUTODETECT" != no ] ; then
		ebegin "Loading hardware drivers"
		# first run we find the controllers 
		# second we find the discs
		find /sys -name modalias | xargs sort -u \
			| xargs modprobe -a  2> /dev/null
		find /sys -name modalias | xargs sort -u \
			| xargs modprobe -a 2> /dev/null
	fi
	eend 0
}

