# init.d/inetd
#!/sbin/runscript

start() {
	ebegin "Starting inetd"
	inetd $INETD_OPTS
	eend $?
}

stop () {
	ebegin "Stopping inetd"
	killall inetd
	eend $?
}
