# init.d/keymap
#!/sbin/runscript


start() {
	[ -z "$KEYMAP" ] && return
	ebegin "Setting keymap"
	zcat "$KEYMAP" | loadkmap
	eend $?
}

stop() {
	return
}
