# init.d/mdev
#!/sbin/runscript

depend() {
	provide dev
	need sysfs mdev-mount
	before checkfs fsck
	keyword novserver
}

start() {
	# check if udev is specified on cmd line
	for i in `cat /proc/cmdline 2>/dev/null`; do
		[ "$i" = "udev" ] && return 0 
	done
	
	ebegin "Starting mdev"
	mkdir -p /dev

	# start mdev
	mount -t tmpfs -o exec,nosuid,mode=0755 mdev /dev
	mknod -m 666 /dev/null c 1 3

	# use mdev for hotplug
	echo "/sbin/mdev" > /proc/sys/kernel/hotplug
	
	# create devices
	mdev -s
	RC=$?

	# create pts file system
	[ -d /dev/pts ] || mkdir -m 755 /dev/pts
	[ -c /dev/ptmx ] || mknod -m 666 /dev/ptmx c 5 2
	mount -t devpts -o gid=5,mode=0620,noexec,nosuid devpts /dev/pts
	[ -d /dev/shm ] || mkdir /dev/shm
	mount -t tmpfs -o nodev,nosuid,noexec shm /dev/shm
	eend $RC
}

stop() {
	ebegin "Stopping mdev"
	echo "" > /proc/sys/kernel/hotplug
	umount /dev/pts && umount /dev/shm && umount /dev
	eend 
}
