# init.d/mdev
#!/sbin/runscript

start() {
	# check if udev is specified on cmd line
	for i in `cat /proc/cmdline`; do
		[ "$i" = "udev" ] && return 0 
	done
	
	ebegin "Starting mdev"
	mkdir -p /dev

	# start mdev
	mount /dev
	mknod -m 666 /dev/null c 1 3

	# use mdev for hotplug
	echo "/sbin/mdev" > /proc/sys/kernel/hotplug
	
	# create devices
	mdev -s
	RC=$?

	# create pts file system
	[ -d /dev/pts ] || mkdir -m 755 /dev/pts
	[ -c /dev/ptmx ] || mknod -m 666 /dev/ptmx c 5 2
	mount /dev/pts
	[ -d /dev/shm ] || mkdir /dev/shm
	mount /dev/shm
	eend $RC
}

stop() {
	ebegin "Stopping mdev"
	echo "" > /proc/sys/kernel/hotplug
	umount /dev/pts && umount /dev/shm && umount /dev
	eend 
}
