# init.d/modutils
#!/sbin/runscript

PATH="/sbin:/bin:/usr/bin:/usr/bin"

start() {
	modules-update
	if [ -f /etc/modules ] ; then
		sed 's/\#.*//g' < /etc/modules |
		while read module args
		do
			ebegin "Loading $module"
			modprobe -q $module $args
			eend $?
		done
	fi
}

