# init.d/procps
#!/sbin/runscript

start() {
	if [ -r /etc/sysctl.conf ] ; then
		ebegin "Setting kernel variables"
		/sbin/sysctl -p /etc/sysctl.conf >/dev/null
		eend $?
	fi
}

