# init.d/rc.local
#!/sbin/runscript

start() {
	if [ -f /etc/rc.local ] ; then
		ebegin "Running local boot script /etc/rc.local"
		/etc/rc.local
		eend $?
	fi
}

stop() {
	return 0
}

