# init.d/swap
#!/sbin/runscript

start() {
	ebegin "Enabling swap"
	swapon -a
	echo -n $( awk '/^SwapTotal/ { print " (" $2 $3 ") " }' /proc/meminfo )
	eend $?
}

stop () {
	ebegin "Disabling swap"
	swapoff -a
	eend $?
}
