# init.d/syslog
#!/sbin/runscript

start() {
	ebegin "Starting system logging"
	klogd ${KLOGD_OPTS}
	syslogd ${SYSLOGD_OPTS}
	eend $?
}

stop () {
	ebegin "Stopping system logging"
	killall klogd
	killall syslogd
	eend $?
}
